# manifest/etc/cron.d/root
5 * * * * /usr/bin/gunzip /reports/*.gz; /usr/bin/dmarcts-report-parser.pl -x --delete -d -r /reports/*.xml 1>>/var/log/nginx/dmarc-reports.log 2>&1

